# infra/systemd/coworker-api.service
[Unit]
Description=MC & S CoWorker v3 API
After=network.target postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=coworker
Group=coworker
WorkingDirectory=/opt/coworker/current/backend
LoadCredentialEncrypted=coworker.env:/opt/coworker/secrets/production.env.age
ExecStartPre=/bin/sh -c 'age -d -i /etc/coworker/age.key -o /run/credentials/coworker-api.service/coworker.env $CREDENTIALS_DIRECTORY/coworker.env'
EnvironmentFile=/run/credentials/coworker-api.service/coworker.env
ExecStart=/opt/coworker/current/.venv/bin/uvicorn coworker.api.main:app \
  --host 127.0.0.1 \
  --port 8001 \
  --workers 2 \
  --log-config /opt/coworker/current/backend/coworker/logging-uvicorn.json
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/coworker /var/log/coworker
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
